# supervisor/roles/worker.conf
[program:worker]
command=php /worker.php
